# etc/rc.local
#!/bin/sh
# 由/etc/rc.local调用
# /etc/rc.local里放入下面三行代码
# if [ -x /root/sys-stuff/etc/rc.local ]; then
#        /root/sys-stuff/etc/rc.local
# fi

CONFIG=`dirname $0`/../etc/config
test -r $CONFIG || exit
. $CONFIG

# 双路由
if [ $DUAL_IP = yes ]; then
    $ROOT/bin/iproute
fi

$ROOT/bin/moodle_init
